# packaging/deb/DEBIAN/postinst
#!/usr/bin/env bash
set -e
systemctl daemon-reload || true
systemctl enable tuxpulse-helper.service || true
